# examples/systemd/git-mirror.timer
[Unit]
Description=Run git-mirror periodically

[Timer]
OnBootSec=5min
OnUnitActiveSec=15min
Unit=git-mirror.service

[Install]
WantedBy=timers.target
